# bench/config/templates/systemd/frappe-bench-frappe-long-worker.service
[Unit]
Description="{{ bench_name }}-frappe-short-worker %I"
PartOf={{ bench_name }}-workers.target

[Service]
User={{ user }}
Group={{ user }}
Restart=always
ExecStart={{ bench_cmd }} worker --queue long
StandardOutput=file:{{ bench_dir }}/logs/worker.log
StandardError=file:{{ bench_dir }}/logs/worker.error.log
WorkingDirectory={{ bench_dir }}
