# bench/config/templates/systemd/frappe-bench-web.target
[Unit]
After=network.target
Wants={{ bench_name }}-frappe-web.service {{ bench_name }}-node-socketio.service

[Install]
WantedBy=multi-user.target
