# bench/config/templates/systemd/frappe-bench-workers.target
[Unit]
After=network.target
Wants={{ worker_target_wants }}

[Install]
WantedBy=multi-user.target
