# bench/config/templates/systemd/frappe-bench-redis.target
[Unit]
After=network.target
Wants={{ bench_name }}-redis-cache.service {{ bench_name }}-redis-queue.service

[Install]
WantedBy=multi-user.target
